# voxtype/voxtype-recordings-sync.timer
[Unit]
Description=Daily sync of voxtype recordings to S3

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
